# hyprland/.config/systemd/user/lid-handler.service
[Unit]
Description=Hyprland Lid Event Handler
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/eric/.local/bin/lid-service-wrapper
Restart=always
RestartSec=3
Environment=HOME=/home/eric
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-1

[Install]
WantedBy=default.target
